# roles/kubernetes/templates/config/kube-sheduler.service
[Unit]
Description=Kubernetes Scheduler
Documentation=https://github.com/GoogleCloudPlatform/kubernetes

[Service]
ExecStart=/opt/kubernetes/bin/kube-scheduler \
  --leader-elect=true \
  --master=http://{{ kubernetes.master }}:8080 \
  --address=127.0.0.1 \
  --v=2
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target